# src.wayland/weston.conf
[program:weston]
command=/usr/bin/weston --socket=wayland-0 
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
